# supervisord.conf
; ============================================================
; CRMEB 微信云托管 Supervisor 配置
; 管理: PHP-FPM · Nginx · 队列消费 · Workerman 长连接 · 定时任务
; ============================================================

[supervisord]
nodaemon=true
logfile=/dev/null
logfile_maxbytes=0
pidfile=/var/run/supervisord.pid
loglevel=info

[unix_http_server]
file=/var/run/supervisor.sock
chmod=0700

[supervisorctl]
serverurl=unix:///var/run/supervisor.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

; ── Redis（最优先启动，队列/缓存/长连接均依赖它）──────────
; 不设密码：Redis 仅监听 127.0.0.1，容器外部无法直接访问
[program:redis]
command=redis-server --bind 127.0.0.1 --port 6379 --daemonize no --loglevel notice
autostart=true
autorestart=true
priority=5
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
stopasgroup=true
killasgroup=true

; ── PHP-FPM（优先启动，Nginx 依赖它）────────────────────
[program:php-fpm]
command=/usr/local/sbin/php-fpm -F -R
autostart=true
autorestart=true
priority=10
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
stopasgroup=true
killasgroup=true

; ── Nginx（HTTP 入口，端口 80）───────────────────────────
[program:nginx]
command=/usr/sbin/nginx -g "daemon off;"
autostart=true
autorestart=true
priority=20
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
stopasgroup=true
killasgroup=true

; ── 队列消费（ThinkPHP Queue）────────────────────────────
; QUEUE_NAME 由 init.sh 从 .env 读取并 export，supervisord 继承该变量
[program:queue]
command=/bin/sh -c "QNAME=$(grep -E '^QUEUE_NAME\\s*=' /var/www/.env 2>/dev/null | head -1 | sed 's/^QUEUE_NAME\\s*=\\s*//' | tr -d ' \\r\\n'); exec php /var/www/think queue:listen --queue \"${QNAME:-crmeb}\""
autostart=true
autorestart=true
priority=30
startsecs=15
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
stopasgroup=true
killasgroup=true

; ── Workerman 客服/通知长连接（端口 40001）───────────────
[program:workerman-chat]
command=php /var/www/think workerman start chat
autostart=true
autorestart=true
priority=30
startsecs=10
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
stopasgroup=true
killasgroup=true

; ── Workerman 后台管理长连接（端口 40002）────────────────
[program:workerman-admin]
command=php /var/www/think workerman start admin
autostart=true
autorestart=true
priority=30
startsecs=10
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
stopasgroup=true
killasgroup=true

; ── 定时任务（ThinkPHP Timer）────────────────────────────
[program:timer]
command=php /var/www/think timer start
autostart=true
autorestart=true
priority=30
startsecs=10
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
stopasgroup=true
killasgroup=true
